# tg-bot.service.example
[Unit]
  Description=BudgetTGBot
 
[Service]
  User=<YOUR_USERNAME>
  Group=<YOUR_GROUP>
  SyslogIdentifier=budget-tg-bot
  SyslogFacility=daemon
  Type=simple
  WorkingDirectory=<PROJECT_DIR>
  Environment=PATH=<PROJECT_DIR>/venv/bin:$PATH
  ExecStart=<PROJECT_DIR>/venv/bin/python3 tg_bot.py
  Restart=always
 
[Install]
  WantedBy=multi-user.target
